# debian/openss7-base.prerm.in
#!/bin/sh
# @configure_input@
# =============================================================================
# 
# @(#) $RCSfile$ $Name$($Revision$) $Date$
#
# -----------------------------------------------------------------------------
#
# Copyright (c) 2008-2011  Monavacon Limited <http://www.monavacon.com/>
# Copyright (c) 2001-2008  OpenSS7 Corporation <http://www.openss7.com/>
# Copyright (c) 1997-2001  Brian F. G. Bidulock <bidulock@openss7.org>
#
# All Rights Reserved.
#
# This program is free software; you can redistribute it and/or modify it under
# the terms of the GNU Affero General Public License as published by the Free
# Software Foundation; version 3 of the License.
#
# This program is distributed in the hope that it will be useful, but WITHOUT
# ANY WARRANTY; without even the implied warranty of MERCHANTABILITY or FITNESS
# FOR A PARTICULAR PURPOSE.  See the GNU Affero General Public License for more
# details.
#
# You should have received a copy of the GNU Affero General Public License along
# with this program.  If not, see <http://www.gnu.org/licenses/>, or write to
# the Free Software Foundation, Inc., 675 Mass Ave, Cambridge, MA 02139, USA.
#
# -----------------------------------------------------------------------------
#
# U.S. GOVERNMENT RESTRICTED RIGHTS.  If you are licensing this Software on
# behalf of the U.S. Government ("Government"), the following provisions apply
# to you.  If the Software is supplied by the Department of Defense ("DoD"), it
# is classified as "Commercial Computer Software" under paragraph 252.227-7014
# of the DoD Supplement to the Federal Acquisition Regulations ("DFARS") (or any
# successor regulations) and the Government is acquiring only the license rights
# granted herein (the license rights customarily provided to non-Government
# users).  If the Software is supplied to any unit or agency of the Government
# other than DoD, it is classified as "Restricted Computer Software" and the
# Government's rights in the Software are defined in paragraph 52.227-19 of the
# Federal Acquisition Regulations ("FAR") (or any successor regulations) or, in
# the cases of NASA, in paragraph 18.52.227-86 of the NASA Supplement to the FAR
# (or any successor regulations).
#
# -----------------------------------------------------------------------------
#
# Commercial licensing and support of this software is available from OpenSS7
# Corporation at a fee.  See http://www.openss7.com/
#
# -----------------------------------------------------------------------------
#
# Last Modified $Date$ by $Author$
#
# =============================================================================
subpkg=base
@openss7_header@
datarootdir=@datarootdir@
# -----------------------------------------------------------------------------
# Auxiliary functions:
# -----------------------------------------------------------------------------
#
# prerm_rmnods() - Remove STREAMS devices generated by STRCONF.  Note that this
#   is a shell script that does not require any STREAMS libraries to be present
#   and so can be called regardless of whether the -lib subpackage is available.
#
prerm_rmnods () { :
    if [ -x "@sbindir@/@STRCONF_DEVICES@" ]; then
	echo "${package}: removing old devices."
	@sbindir@/@STRCONF_DEVICES@ --remove || :
    fi
    return 0
}
#
# prerm_initscripts() - Before removing or upgrading a package, stop the init
#   scripts.  This is fine for 2.4 module systems, but it is a little too
#   disruptive for weak-module updates in the 2.6 module system.  Because this
#   is too disruptive for weak-updated systems, we disable the debhelper
#   scriplets and only stop init scripts when not performing weak updates, and
#   even then, only when the updates are for the running kernel.
#
prerm_initscripts () {
    local reverse= list="${initfiles}" i
    if [ "@kversion@" = "`uname -r`" ]; then
	for i in $list; do
	    reverse="$i${reverse:+' '}$reverse"
	done
	for i in $reverse; do
	    if [ -f "@initrddir@/$i" ] ; then
		[ -x "`which invoke-rc.d 2>/dev/null`" ] && invoke-rc.d $i stop && continue || true
		[ -x "`which service     2>/dev/null`" ] && service     $i stop && continue || true
		[ -x "@initrddir@/$i"                  ] && @initrddir@/$i stop && continue || true
	    fi
	done
    fi
    :
}
prerm_manrefs () { :
}
prerm_manmacros () { :
}
#
# prerm_base_remove() - Remove the -base subpackage completely.  No upgrade is
#   in processes.
prerm_base_remove () { :
    # We are being removed for good.
    prerm_rmnods
    @WITH_KO_MODULES_FALSE@prerm_initscripts
}
prerm_base_upgrade () { :
    # We are being removed because we are being replaced with a newer version.
    @WITH_KO_MODULES_FALSE@prerm_initscripts
    prerm_manrefs
    prerm_manmacros
}
# -----------------------------------------------------------------------------
# Primary functions:
# -----------------------------------------------------------------------------
# xxx-prerm remove
remove_prerm () { :
    prerm_base_remove
}
# con-prerm remove 'in-favour' package new-version
conflict_prerm () { :
    prerm_base_remove
}
# old-prerm upgrade new-version
upgrade_prerm () { :
    prerm_base_upgrade
}
# old-prerm upgrade new-version ***FAILED***
# new-prerm failed-upgrade old-version
failed_upgrade_prerm () { :
    prerm_base_upgrade
}
# dec-prerm deconfigure 'in-favour' package-being-installed version ['removing' conflicting-package version]
deconfigure_prerm () { :
    prerm_base_upgrade
}
@mscript_prerm@
@mscript_debhelper@
# =============================================================================
#
# $Log$
# =============================================================================
# vim: ft=sh sw=4 noet nocin nosi formatoptions+=tcqlorn
