# debian/openss7-base.postrm.in
#!/bin/bash
# @configure_input@
# =============================================================================
# 
# @(#) $RCSfile: openss7-base.postrm.in,v $ $Name:  $($Revision: 1.1.2.2 $) $Date: 2011-05-15 12:01:18 $
#
# -----------------------------------------------------------------------------
#
# Copyright (c) 2008-2011  Monavacon Limited <http://www.monavacon.com/>
# Copyright (c) 2001-2008  OpenSS7 Corporation <http://www.openss7.com/>
# Copyright (c) 1997-2001  Brian F. G. Bidulock <bidulock@openss7.org>
#
# All Rights Reserved.
#
# This program is free software; you can redistribute it and/or modify it under
# the terms of the GNU Affero General Public License as published by the Free
# Software Foundation; version 3 of the License.
#
# This program is distributed in the hope that it will be useful, but WITHOUT
# ANY WARRANTY; without even the implied warranty of MERCHANTABILITY or FITNESS
# FOR A PARTICULAR PURPOSE.  See the GNU Affero General Public License for more
# details.
#
# You should have received a copy of the GNU Affero General Public License along
# with this program.  If not, see <http://www.gnu.org/licenses/>, or write to
# the Free Software Foundation, Inc., 675 Mass Ave, Cambridge, MA 02139, USA.
#
# -----------------------------------------------------------------------------
#
# U.S. GOVERNMENT RESTRICTED RIGHTS.  If you are licensing this Software on
# behalf of the U.S. Government ("Government"), the following provisions apply
# to you.  If the Software is supplied by the Department of Defense ("DoD"), it
# is classified as "Commercial Computer Software" under paragraph 252.227-7014
# of the DoD Supplement to the Federal Acquisition Regulations ("DFARS") (or any
# successor regulations) and the Government is acquiring only the license rights
# granted herein (the license rights customarily provided to non-Government
# users).  If the Software is supplied to any unit or agency of the Government
# other than DoD, it is classified as "Restricted Computer Software" and the
# Government's rights in the Software are defined in paragraph 52.227-19 of the
# Federal Acquisition Regulations ("FAR") (or any successor regulations) or, in
# the cases of NASA, in paragraph 18.52.227-86 of the NASA Supplement to the FAR
# (or any successor regulations).
#
# -----------------------------------------------------------------------------
#
# Commercial licensing and support of this software is available from OpenSS7
# Corporation at a fee.  See http://www.openss7.com/
#
# -----------------------------------------------------------------------------
#
# Last Modified $Date: 2011-05-15 12:01:18 $ by $Author: brian $
#
# =============================================================================
subpkg=base
@openss7_header@
datarootdir=@datarootdir@
# -----------------------------------------------------------------------------
# Auxiliary functions:
# -----------------------------------------------------------------------------
#
# postrm_oldmods() - As the subpackage is being purged, remove the line from the modules.conf file
#   that points at the /etc/modutils/openss7 file that in turn includes the /lib/modules/@kversion@-
#   /modules.openss7 file when it exists.  Note that the files themselves referenced from the
#   modules.conf files are removed as a normal part of the course of removing the files in the
#   -kernel sub-packages.
#
postrm_oldmods () {
    local conf_old conf_tmp conf_new
    [ -z '@WITH_KO_MODULES_FALSE@' ] || return 0
    conf_old=${sysconfdir}/modules.conf
    conf_tmp=${tmppath:-/var/tmp}/modules.conf.tmp.$$
    conf_new=${tmppath:-/var/tmp}/modules.conf.new.$$
    if [ -f $conf_old ] ; then
        cp -f -- $conf_old $conf_new
        if grep '^include[[:space:]].*'${modutildir}/${PACKAGE_LCNAME} $conf_new >/dev/null 2>&1; then
            sed -e '\|^include[[:space:]].*'${modutildir}/${PACKAGE_LCNAME}'|d' $conf_new > $conf_tmp
            mv -f -- $conf_tmp $conf_new
            cp -f -b --suffix=.dpkg-old -- $conf_new $conf_old
        fi
        rm -f -- $conf_new || :
    fi
}
#
# postrm_initscripts() - When purging the package, remove the symbolic links associated with the
#   init scripts before they are purged.  Note that the debhelper scriptlets are too disruptive in
#   some stages and are therefore disabled.  Because some time could have elapsed between package
#   removal and package purge, unlike debhelper, we also stop any running initscripts.
#
postrm_initscripts () {
    local reverse list i
    reverse=
    list="${initfiles}"
    for i in $list; do
	reverse="$i${reverse:+ }$reverse"
    done
    for i in $reverse; do
	if [ -f ${initrddir}/$i ]; then
	    [ -x "`which invoke-rc.d 2>/dev/null`" ] && { invoke-rc.d  $i stop && continue || :; }
	    [ -x "`which service     2>/dev/null`" ] && { service      $i stop && continue || :; }
	    [ -x ${initrddir}/$i                   ] && { ${initrddir}/$i stop && continue || :; }
	fi
    done
    for i in $reverse; do
	if [ -f ${initrddir}/$i ]; then
	    [ -x /sbin/insserv             ] && { /sbin/insserv -r          $i        2>&1 && continue || :; }
	    [ -x /sbin/chkconfig           ] && { /sbin/chkconfig --del     $i        2>&1 && continue || :; }
	    [ -x /usr/lib/lsb/remove_initd ] && { /usr/lib/lsb/remove_initd $i        2>&1 && continue || :; }
	    [ -x /usr/sbin/update-rc.d     ] && { /usr/sbin/update-rc.d     $i remove 2>&1 && continue || :; }
	fi
    done
}
#
# postrm_manpages() - When removing the package, update the aprospros database.  This handled by
#   debhelper scriptlets or triggers and is therefore short-circuited.
#
postrm_manpages () { :; return
    if   [ -x "`which mandb      2>/dev/null`" ]; then
	LANGUAGE= LANG= mandb   -q -p ${mandir} || :
    elif [ -x "`which makewhatis 2>/dev/null`" ]; then
	LANGUAGE= LANG= makewhatis -u ${mandir} || :
    fi
}
#
# postrm_purge() - The subpackage is being purged, rebuild databases and adjust other files in
#   preparation for removal of configuration files.
#
postrm_purge () { :
    postrm_oldmods || :
    postrm_manpages || :
}
# -----------------------------------------------------------------------------
# Primary functions:
# -----------------------------------------------------------------------------
# remove-postrm remove
remove_postrm () { :
}
# remove-postrm purge
purge_postrm () { :
    postrm_purge || :
}
# old-postrm upgrade new-version
upgrade_postrm () { :
}
# old-postrm upgrade new-version ...FAILED
# new-postrm failed-upgrade old-version
failed_upgrade_postrm () { :
}
# new-preinst install [old-version] ...FAILED
# new-postrm abort-install [old-version]
abort_install_postrm () { :
}
# new-preinst upgrade old-version ...FAILED
# new-postrm abort-upgrade old-version
abort_upgrade_postrm () { :
}
# dis-postrm disappear overwriter overwriter-version
disappear_postrm () { :
}
@mscript_postrm@
@mscript_debhelper@
# =============================================================================
#
# $Log: openss7-base.postrm.in,v $
# Revision 1.1.2.2  2011-05-15 12:01:18  brian
# - many build system refinements
#
# Revision 1.1.2.1  2011-05-10 13:45:28  brian
# - weak modules workup
#
# =============================================================================
# vim: ft=sh sw=4 noet nocin nosi formatoptions+=tcqlorn
