# src/util/systemd/streams-shutdown.target
[Unit]
Description=Shutdown of STREAMS Subsystem
Documentation=man:streams(4)
