# src/util/systemd/strerr.service.in
[Unit]
Description=STREAMS error logger
Documentation=man:strerr(8)
ConditionPathIsMountPoint=/dev/streams
Conflicts=streams-shutdown.target
Before=streams-shutdown.target
After=streams.service

[Service]
Type=forking
EnvironmentFile=-/etc/sysconfig/strerr
EnvironmentFile=-/etc/default/strerr
EnvironmentFile=-/etc/conf.d/strerr
PIDFile=/var/run/strerr.pid
ExecStartPre=@sysbindir@/mkdir -p -- /var/log/streams
ExecStart=@sbindir@/strerr $STRERR_OPTIONS
ExecReload=@sysbindir@/kill -HUP $MAINPID

[Install]
WantedBy=streams.target
